# Keys Remap/Fedora (universal)/files/etc/systemd/system/set-keys-mapping.service
[Unit]
Description=Setting keys mapping

[Service]
Type=oneshot
ExecStart=/usr/bin/set-keys-mapping

[Install]
WantedBy=multi-user.target
